# roles/void-mirror/files/void-mirror.cron
*/5 * * * *  reposync rsync -avzurk --delete repo3.voidlinux.eu::repo/ /mnt/mirrordata/void >/dev/null 2>&1
